# meta-intel-iot-middleware/recipes-devtools/wyliodrin-server/files/wyliodrin-server.service
[Unit]
Description=Wyliodrin server
After=redis
ConditionFileNotEmpty=/media/card/wyliodrin.json

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/wyliodrind
Restart=always

[Install]
WantedBy=multi-user.target
