# meta-ostro-bsp/meta-edison-bsp/recipes-connectivity/bluetooth-rfkill-event/files/bluetooth-rfkill-event.service
[Unit]
Description=Bluetooth rf kill event daemon

[Service]
Type=simple
ExecStart=/usr/sbin/bluetooth_rfkill_event
Restart=on-failure

[Install]
WantedBy=multi-user.target
