# meta-appfw/recipes-appfw/example-session/files/example-corp.service
[Unit]
Description=Example session
Requires=tlm.service iot-launch.socket
After=tlm.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/tlm-client -k -l --seat seat_placeholder --username example-corp --password automatic_login
