# meta-soletta/recipes-soletta/dev-app/soletta-dev-app/soletta-dev-app-avahi-discover.service
[Unit]
Description=Autogenerate hostname with MAC ADDRESS

[Service]
Type=oneshot
ExecStart=/opt/soletta-dev-app/scripts/soletta-dev-app-mac.sh

[Install]
WantedBy=multi-user.target

