# meta-ostro/recipes-connectivity/caninitscripts/files/can-init.service
[Unit]
Description=CAN startup service
After=network.target
Type=oneshot

[Service]
ExecStart=/usr/bin/start-can

[Install]
WantedBy=network.target
