# meta-swupd/recipes-core/swupd-client/oe-swupd-helpers/ldconfig-trigger.service
[Unit]
Description=Rebuild Dynamic Linker Cache
Documentation=man:ldconfig(8)
BindsTo=update-triggers.target

[Service]
Type=oneshot
ExecStart=/sbin/ldconfig -X
