# meta-iot-web/recipes-web/iot-rest-api-server/files/iot-rest-api-server.socket
[Unit]
Description=Iot REST API server socket

[Socket]
ListenStream=8000
Accept=false

[Install]
WantedBy=sockets.target
