# meta-iot-web/recipes-smarthome/smarthome-ocf-servers/files/smarthome-rgbled.service
[Unit]
Description=SmartHome RGB LED sensor startup service
After=network.target

[Service]
ExecStart=/usr/bin/node /opt/smarthome-ocf-servers/smarthome-rgbled.js
Environment='NODE_PATH=/usr/lib/node_modules/'
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
