# meta-iot-web/recipes-smarthome/smarthome-gateway/files/smarthome-gateway.socket
[Unit]
Description=SmartHome Gateway socket

[Socket]
ListenStream=8080
Accept=false

[Install]
WantedBy=sockets.target
